# flask_bbs.conf
[program:flask_bbs]
command=/usr/local/bin/gunicorn wsgi --bind 0.0.0.0:2000
directory=/var/flask_bbs
autostart=true
autorestart=true
